# ansible/roles/rocketmq/files/opt/app/current/conf/systemd/rocketmq-namesrv.service
[Unit]
Description=RocketMQ Nameserver
After=network.target

[Service]
Type=simple
User=rocketmq
Group=svc
LimitCORE=infinity
LimitNOFILE=infinity
EnvironmentFile=/opt/app/current/bin/envs/java.env
EnvironmentFile=/opt/app/current/conf/namesrv/namesrv.env
ExecStart=/opt/rocketmq/current/bin/mqnamesrv
ExecStop=/opt/rocketmq/current/bin/mqshutdown namesrv
Restart=always

[Install]
WantedBy=multi-user.target
